# deb/DEBIAN/prerm
#!/bin/sh

set -e
echo "Stopping routerest service"
systemctl stop routerest
systemctl disable routerest
